# base-system/etc/udev/rules.d/09-watchdog.rules
DEVPATH=="/bus/i2c/drivers/rtc-ds1374", RUN+="/sbin/modprobe ds1374_wdt"
